# systemd/cunctator.service.in
@GENERATED_SCRIPT_FILE_WARNING@

[Unit]
Description=Cunctator Service
After=network.target remote-fs.target nss-lookup.target
Documentation=man:cuncd(8)

[Service]
Type=simple
ExecStart=@prefix@/sbin/cuncd
TimeoutStopSec=10
PrivateTmp=true
Restart=always

[Install]
WantedBy=multi-user.target
